# systemd/silverblue-recipe@.service.example
[Unit]
Description=Regenerate the silverblue-recipe OS-layer backup for %i
# Uncomment to order the capture before a btrbk backup:
#Before=btrbk.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/silverblue-recipe capture --user %i

[Install]
WantedBy=multi-user.target
